# filen-sync/filen-sync.timer
[Unit]
Description=Run Filen Sync

[Timer]
OnStartupSec=1min
OnUnitInactiveSec=5min

[Install]
WantedBy=timers.target
